# sbin/backup_warehouse.crontab
0 1 * * * /bin/bash -l -c %APP_HOME%/sbin/backup_warehouse.sh
