# sbin/database_backup.crontab
0 21 * * * /bin/bash -l -c %APP_HOME%/sbin/database_backup.sh
